# online-kalkulator-server/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
source /usr/lib/bashio/bashio.sh

exec php-fpm -F

