# online-kalkulator-server/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
source /usr/lib/bashio/bashio.sh

echo Listing PHP versions
ls /usr/sbin | grep php

# Running PHP
exec php-fpm82 -F

