# online-kalkulator-server/rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash
exec mysqld --datadir=/data/mysql --user=mysql